# systemd/user/bootc-bootstrap.service
[Unit]
Description=Apply user bootstrap state (Flatpak + GNOME extensions)
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bootc-bootstrap apply

[Install]
WantedBy=default.target
